# scripts/systemd/finance-app-updater.path
[Unit]
Description=Watch for Finance App update trigger
After=docker.service

[Path]
PathExists=/opt/finance-app/update-data/.update-requested

[Install]
WantedBy=multi-user.target
