# scripts/systemd/finance-app-updater.path
[Unit]
Description=Watch for Finance App update trigger
After=docker.service

[Path]
PathExists=/opt/finance-app/.update-requested
MakeDirectory=yes

[Install]
WantedBy=multi-user.target
